# build/service/mysql/run
#!/bin/sh
if [ ! -f /mount/database/ibdata1 ]; then
    mysql_install_db --datadir=/mount/database/
fi
chmod -R a+wX /mount/log
exec /usr/sbin/mysqld